# 99-usb-block.rules
ACTION=="add", SUBSYSTEM=="usb", ENV{DEVTYPE}=="usb_interface", \
ATTRS{bInterfaceClass}=="08", RUN+="/usr/local/bin/usb_blocker.sh $env{DEVPATH}"

